# crates/catalyrst-camera-reel/deploy/catalyrst-camera-reel.service
[Unit]
Description=catalyrst-camera-reel (Rust port of camera-reel-service.decentraland.org, port 5143)
After=network.target postgresql.service catalyrst-places.service
StartLimitIntervalSec=0
Wants=postgresql.service

[Service]
Type=exec
WorkingDirectory=<WORKSPACE>
EnvironmentFile=<WORKSPACE>/env/catalyrst-camera-reel.env
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do pg_isready -h <HOST> -p 5433 -U <DB_USER> -q && exit 0; sleep 1; done; exit 1'
ExecStart=<WORKSPACE>/bin/catalyrst-camera-reel
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
