# crates/catalyrst-profile-images/deploy/catalyrst-profile-images.service
[Unit]
Description=catalyrst-profile-images (Rust port of profile-images.decentraland.org — avatar face/body thumbnails, port 5152)
After=network.target catalyrst.service
StartLimitIntervalSec=0

[Service]
Type=exec
WorkingDirectory=<WORKSPACE>
EnvironmentFile=<WORKSPACE>/env/catalyrst-profile-images.env
ExecStartPre=/bin/sh -c 'mkdir -p <DATA_DIR>/profile-images/.render-tmp'
ExecStart=<WORKSPACE>/bin/catalyrst-profile-images
LimitCORE=0
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
