# crates/catalyrst-price/catalyrst-price.service
[Unit]
Description=catalyrst-price (serving layer for api.coingecko.com price-feed, port 5144)
After=network.target postgresql.service
StartLimitIntervalSec=0
Wants=postgresql.service

[Service]
Type=exec
WorkingDirectory=<WORKSPACE>
EnvironmentFile=<WORKSPACE>/env/catalyrst-price.env
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do pg_isready -h <HOST> -p 5433 -U <DB_USER> -q && exit 0; sleep 1; done; exit 1'
ExecStart=<WORKSPACE>/bin/catalyrst-price
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
